# server/pxf-service/src/templates/conf/pxf.service.d/local.conf
[Service]
Environment=PXF_BASE=@PXF_BASE@
PIDFile=@PXF_BASE@/run/pxf-app.pid